# scripts/systemd/route53_cdk_dyndns.service
[Unit]
Description=Route53 Dyndns Service

[Service]
Type=oneshot
ExecStart=/home/riccardomc/route53-cdk-dyndns/scripts/run.sh
