# scripts/systemd/route53_cdk_dyndns.timer
[Unit]
Description=Route53 Dyndns Timer

[Timer]
OnUnitActiveSec=30min
OnBootSec=30min

[Install]
WantedBy=timers.target
